# packaging/debian/DEBIAN/prerm
#!/bin/bash
set -e

# Pre-removal script for family-policy

# Stop and disable service if running
if command -v systemctl >/dev/null 2>&1; then
    if systemctl is-active --quiet family-policy-agent; then
        systemctl stop family-policy-agent || true
    fi
    if systemctl is-enabled --quiet family-policy-agent 2>/dev/null; then
        systemctl disable family-policy-agent || true
    fi
fi

exit 0
